# com.biodisc.autonomous.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.biodisc.autonomous</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Library/Frameworks/Python.framework/Versions/3.14/bin/python3</string>
        <string>/Users/gjw255/astrodata/SWARM/BIODISC/start_biodisc_autonomous.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/gjw255/astrodata/SWARM/BIODISC</string>

    <key>StandardOutPath</key>
    <string>/Users/gjw255/astrodata/SWARM/BIODISC/biodisc_autonomous_stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/gjw255/astrodata/SWARM/BIODISC/biodisc_autonomous_stderr.log</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>RestartInterval</key>
    <integer>300</integer>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONPATH</key>
        <string>/Users/gjw255/astrodata/SWARM/BIODISC</string>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/Library/Frameworks/Python.framework/Versions/3.14/bin</string>
    </dict>

    <key>ProcessType</key>
    <string>Background</string>

    <key>Nice</key>
    <integer>10</integer>
</dict>
</plist>
